# launchd/macmini/com.myscope.layer3-index.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Mac mini 定时任务：第三层 RSS 索引
  FreshRSS → Meilisearch hubble_radius
  每天 02:30 运行
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.myscope.layer3-index</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>MYSCOPE_DIR/scripts/layer3_index.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>MYSCOPE_DIR</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>MYSCOPE_DIR/logs/layer3_index.log</string>
  <key>StandardErrorPath</key>
  <string>MYSCOPE_DIR/logs/layer3_index_err.log</string>
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
