# launchd/macmini/com.myscope.layer1-flomo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Mac mini 定时任务：第一层 flomo 采集
  flomo 网页 → DeepSeek 切片 → Meilisearch memory_chunks
  每天 19:10 运行；失败后在当晚自动补偿重试（DeepSeek 夜间窗口）
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.myscope.layer1-flomo</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>MYSCOPE_DIR/scripts/layer1_flomo.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>MYSCOPE_DIR</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>19</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>19</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>19</integer>
      <key>Minute</key>
      <integer>40</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>20</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
  </array>
  <key>StandardOutPath</key>
  <string>MYSCOPE_DIR/logs/layer1_flomo.log</string>
  <key>StandardErrorPath</key>
  <string>MYSCOPE_DIR/logs/layer1_flomo_err.log</string>
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
