# launchd/macmini/com.myscope.hippocampus-formation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Mac mini 定时任务：海马体 Formation
  Codex / Hermes / Clacky 会话日志 → 过滤 → Anda Formation API
  每天 06:00 运行
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.myscope.hippocampus-formation</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>MYSCOPE_DIR/scripts/hippocampus_formation.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>MYSCOPE_DIR</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>MYSCOPE_DIR/logs/hippocampus_formation.log</string>
  <key>StandardErrorPath</key>
  <string>MYSCOPE_DIR/logs/hippocampus_formation_err.log</string>
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
