# launchd/macmini/com.myscope.layer2-wiki.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Mac mini 定时任务：第二层 LLM Wiki
  Meilisearch memory_chunks → DeepSeek 归纳 → Meilisearch wiki_entries
  每天 05:30 运行（layer1 跑完后）
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.myscope.layer2-wiki</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>MYSCOPE_DIR/scripts/layer2_wiki.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>MYSCOPE_DIR</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>5</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>MYSCOPE_DIR/logs/layer2_wiki.log</string>
  <key>StandardErrorPath</key>
  <string>MYSCOPE_DIR/logs/layer2_wiki_err.log</string>
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
