# launchd/com.hermes.layer1-rag.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  第一层：RAG 切片
  Notion 最近页面 → DeepSeek 切片 → Notion RAG DB
  每天 05:00 运行（layer3 跑完后）
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hermes.layer1-rag</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>HERMES_DIR/scripts/layer1_rag.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>HERMES_DIR</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>5</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>HERMES_DIR/logs/layer1.log</string>
  <key>StandardErrorPath</key>
  <string>HERMES_DIR/logs/layer1_err.log</string>

  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
