# launchd/com.hermes.hippocampus-formation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  海马体记忆形成
  Claude 日志 → DeepSeek 提取 → Notion 海马体数据库
  每天 06:00 运行
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hermes.hippocampus-formation</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>HERMES_DIR/scripts/hippocampus_formation.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>HERMES_DIR</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>HERMES_DIR/logs/hippocampus_formation.log</string>
  <key>StandardErrorPath</key>
  <string>HERMES_DIR/logs/hippocampus_formation_err.log</string>

  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
